# imraa/imraa.service
[Unit]
Description=imraa

[Service]
Type=oneshot
ExecStart=/usr/bin/imraa

[Install]
WantedBy = multi-user.target
