# imraa/imraa.service
[Unit]
Description=imraa

[Service]
Type=oneshot
ExecStart=/usr/bin/imraa -a

[Install]
WantedBy=multi-user.target
